# deploy/filest.rc
#!/bin/sh

# PROVIDE: filest
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="filest"
rcvar="filest_enable"

load_rc_config $name

: ${filest_enable:="NO"}
: ${filest_user:="www"}
: ${filest_env_file:="/usr/local/etc/filest/env"}
: ${filest_log:="/var/log/filest.log"}

pidfile="/var/run/${name}/${name}.pid"
command="/usr/local/bin/filest"
command_interpreter=""

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

filest_prestart()
{
    if [ ! -f "${filest_env_file}" ]; then
        echo "${filest_env_file} not found"
        return 1
    fi
    mkdir -p /var/run/${name}
    chown ${filest_user} /var/run/${name}
    touch ${filest_log}
    chown ${filest_user} ${filest_log}
}

filest_start()
{
    filest_prestart || return 1
    echo "Starting ${name}."
    /usr/sbin/daemon -f -P ${pidfile} -o ${filest_log} \
        /usr/bin/su -m ${filest_user} -c ". ${filest_env_file} && exec ${command}"
}

filest_stop()
{
    if [ -f "${pidfile}" ]; then
        echo "Stopping ${name}."
        kill $(cat ${pidfile}) 2>/dev/null
        rm -f ${pidfile}
    else
        echo "${name} is not running."
    fi
}

filest_status()
{
    if [ -f "${pidfile}" ] && kill -0 $(cat ${pidfile}) 2>/dev/null; then
        echo "${name} is running as pid $(cat ${pidfile})."
    else
        echo "${name} is not running."
        return 1
    fi
}

run_rc_command "$1"
